# deploy/tiptag-unified-incremental.timer
[Unit]
Description=Run TipTag Robinhood unified incremental sink every minute

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=1s
RandomizedDelaySec=0
Persistent=true
Unit=tiptag-unified-incremental.service

[Install]
WantedBy=timers.target
